# src/fty-discovery.service.in
[Unit]
Description=fty-discovery service
After=malamute.service network.target
Requires=malamute.service network.target
Before=fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=discovery-monitoring-daemon
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/lib/nut:"
Environment="prefix=@prefix@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
Environment="NUT_ALTPIDPATH=\"\""
Environment="NUT_STATEPATH=/tmp"
ExecStart=@prefix@/bin/fty-discovery --agent
Restart=always

[Install]
WantedBy=bios.target
