# src/fty-discovery.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=fty-discovery service
After=malamute.service network.target bios-db-init.service 
Requires=malamute.service network.target bios-db-init.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/lib/nut:"
Environment="prefix=@prefix@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
ExecStart=@prefix@/bin/fty-discovery @sysconfdir@/@PACKAGE@/fty-discovery.cfg
Restart=always

[Install]
WantedBy=bios.target
